# src/installation/pytrain_client.service.template
[Unit]
Description=PyTrain Client
After=multi-user.target

[Service]
Type=exec
User=___USER___
Group=___USER___
ExecStart=___HOME___/pytrain_client.bash
Nice=-5

[Install]
WantedBy=multi-user.target